# crontab.cron
MAILTO="clara@clara-tan.home"
0,30 * * * *    /bin/sh /opt/synctorrents/synctorrents.sh >> /var/log/synctorrents/synctorrents_cron.log 2&>1
