# tasks_custom.cron
##This is for /etc/cron.d/task_custom
MAILTO="clara@clara-tan.home"
5 8 * * 0       root    /bin/sh /opt/ubuntu-server-backup/perform.sh                    #Sunday 08:05am maintenance see /home/clara/tools/sunday.sh
20 8 15 * *     root    /bin/sh /home/clara/tools/diskstats.sh                    #15th Disk stats see /home/clara/tools/diskstats
